# systemd_service/map_app.service
[Unit]
Description=IOT Scooter Map Application
After=network.target map_backend.service
Wants=map_backend.service

[Service]
Type=simple
User=iot
Group=iot
WorkingDirectory=/home/iot/map_app_src

# Environment configuration
Environment="PATH=/home/iot/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="LD_LIBRARY_PATH=/usr/lib/aarch64-linux-gnu:/usr/local/lib/aarch64-linux-gnu"
Environment="FLUTTER_DISABLE_CRASH_REPORTING=true"

# Give backend time to start and initialize
ExecStartPre=/bin/sleep 3

# Start with flutter-pi on DRM/KMS (no desktop session required)
ExecStart=/home/iot/map_app/flutter-pi --release /home/iot/map_app_src/build/flutter_assets

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitInterval=120
StartLimitBurst=3

# Timeout settings
TimeoutStartSec=90
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
